# reboot-cycle-script/reboot-cycle.service
[Unit]
Description=Automated Reboot Cycle Counter
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/reboot-cycle/reboot_cycle.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
